# templates/systemd.service.erb
## managed by Puppet ##
[Unit]
Description=<%= @name %>
Wants=basic.target
After=basic.target network.target

[Service]
Type=simple
ExecStart=<%= @real_command %>
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
Restart=always
RestartSec=42s

[Install]
WantedBy=multi-user.target
